# docker/app/supervisord.conf
[supervisord]
nodaemon=true
logfile = /sources/var/log/supervisord.log
pidfile = /sources/var/log/supervisord.pid

[program:symfony]
command=/bin/bash -c "symfony composer install && symfony serve --no-tls --port 8000"
directory=/sources
user=dev
stdout_logfile=/sources/var/log/symfony.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/sources/var/log/symfony-err.log
stdout_logfile_maxbytes=1MB

[program:yarn]
command=/bin/bash -c "yarn && yarn watch"
directory=/sources
user=dev
stdout_logfile=/sources/var/log/yarn.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/sources/var/log/yarn-err.log
stdout_logfile_maxbytes=1MB
